# debian/postrm
#!/bin/bash
# Post-removal cleanup

set -e

echo "Linux Hello - Cleanup"
echo "===================="

# Stop and disable daemon (silent)
systemctl --user stop hello-daemon.service 2>/dev/null || true
systemctl --user disable hello-daemon.service 2>/dev/null || true

# Restaurer le MainBlock.qml original de plasma-workspace
PLASMA_QML="/usr/share/plasma/shells/org.kde.plasma.desktop/contents/lockscreen/MainBlock.qml"
if dpkg-divert --list "$PLASMA_QML" 2>/dev/null | grep -q linux-hello; then
    rm -f "$PLASMA_QML"
    dpkg-divert --package linux-hello --rename --remove "$PLASMA_QML" 2>/dev/null || true
    echo "✓ Lock screen QML original restauré"
fi

# Restore PAM configuration from most recent backup (automatic, non-interactive)
if compgen -G "/etc/pam.d/sudo.pre-linuxhello-*" > /dev/null; then
    echo "Restoring PAM configuration from backup..."
    LATEST_SUDO=$(find /etc/pam.d -name "sudo.pre-linuxhello-*" -type f -printf '%T@ %p\n' 2>/dev/null | sort -rn | head -1 | cut -d' ' -f2-)
    if [ -n "$LATEST_SUDO" ] && [ -f "$LATEST_SUDO" ]; then
        cp "$LATEST_SUDO" /etc/pam.d/sudo
        echo "✓ Restored /etc/pam.d/sudo"
    fi
fi

echo ""
echo "Linux Hello uninstalled"
echo ""

exit 0
