# debian/postinst
#!/bin/bash
set -e

VENV="/opt/linux-hello/venv"
PYTHON="/usr/bin/python3"

. /usr/share/debconf/confmodule
export TEXTDOMAIN=linux-hello
export TEXTDOMAINDIR=/usr/share/locale

echo "$(gettext "[linux-hello] Configuration du paquet…")"

# 1. Créer les dossiers de données (config partagée uniquement)
mkdir -p /var/lib/linux-hello
chmod 755 /var/lib/linux-hello

# 2. Créer ou réparer le venv
echo "[linux-hello] Vérification du venv…"
/usr/lib/linux-hello/repair_venv.sh || true

# 2b. Fixer les permissions du venv pour que linux-hello puisse l'utiliser
echo "[linux-hello] Ajustement des permissions du venv…"
if id -u linux-hello >/dev/null 2>&1; then
    chown -R linux-hello:linux-hello "$VENV" || true
fi
chmod -R u+rw,g+r,o+r "$VENV" || true

# 3. Installer onnxruntime si absent
if ! "$VENV/bin/python3" -c "import onnxruntime" 2>/dev/null; then
    echo "[linux-hello] Installation de onnxruntime via pip…"
    "$VENV/bin/pip" install onnxruntime
fi

# 4. Modifier le shebang du daemon pour utiliser le venv
sed -i '1c #!/opt/linux-hello/venv/bin/python3' /usr/bin/linux-hello-daemon
chmod +x /usr/bin/linux-hello-daemon

# 5. Modifier le shebang de la CLI pour utiliser le venv
sed -i '1c #!/opt/linux-hello/venv/bin/python3' /usr/bin/hello
chmod +x /usr/bin/hello

# 6. Activer le service systemd si disponible
if command -v systemctl >/dev/null 2>&1; then
    echo "[linux-hello] Activation du service systemd…"
    systemctl daemon-reload || true
    systemctl enable --now linux-hello.service || true
else
    echo "[linux-hello] systemd non disponible (WSL2 ?), service non activé."
fi

# 7. Installer et configurer le module PAM
echo "[linux-hello] Installation du module PAM…"
chmod 644 /lib/x86_64-linux-gnu/security/pam_linux_hello.so 2>/dev/null || true
chmod 755 /usr/lib/linux-hello/linux-hello-pam-confirm.py 2>/dev/null || true
chmod 755 /usr/lib/linux-hello/enable-pam-face-recognition.sh 2>/dev/null || true

# Create PAM config directory
mkdir -p /run/linux-hello-pam
chmod 1777 /run/linux-hello-pam

echo ""
echo "[linux-hello] Module PAM installe!"
echo ""
echo "Pour ACTIVER la reconnaissance faciale pour un service :"
echo "    sudo /usr/lib/linux-hello/enable-pam-face-recognition.sh sudo"
echo ""
echo "Services disponibles : sudo, login, sshd, gdm-password"
echo ""
echo "ATTENTION : Testez avec 'pamtester' avant de modifier /etc/pam.d/sudo !"
echo "    sudo pamtester sudo \$USER authenticate"
echo ""