# debian/postinst
#!/bin/bash
# Postinstallation script for linux-hello

set -e

echo "Linux Hello - Postinstallation Configuration"
echo "=============================================="
echo ""

# Create storage directory
STORAGE_DIR="${HOME}/.local/share/linux-hello"
mkdir -p "$STORAGE_DIR" 2>/dev/null || true
chmod 700 "$STORAGE_DIR" 2>/dev/null || true

echo "✓ Daemon service installed"

# Create /run/hello-pam/ immediately (will be recreated at boot by systemd-tmpfiles)
systemd-tmpfiles --create /usr/lib/tmpfiles.d/linux-hello.conf 2>/dev/null || \
    mkdir -p /run/hello-pam && chmod 1777 /run/hello-pam || true
echo "✓ Socket directory /run/hello-pam/ ready"

# Install the KDE lock screen QML overlay via dpkg-divert.
# dpkg-divert moves the plasma-workspace file (renamed .distrib) and puts ours in its place.
# On a plasma-workspace update, dpkg sees the diversion and does not overwrite our file.
PLASMA_QML="/usr/share/plasma/shells/org.kde.plasma.desktop/contents/lockscreen/MainBlock.qml"
LINUX_HELLO_QML="/usr/share/linux-hello/lockscreen/MainBlock.qml"
if [ -f "$LINUX_HELLO_QML" ] && [ -f "$PLASMA_QML" ]; then
    dpkg-divert --package linux-hello --rename --add "$PLASMA_QML" 2>/dev/null || true
    cp "$LINUX_HELLO_QML" "$PLASMA_QML"
    echo "✓ Lock screen QML overlay installed"
fi

# Enable and start the daemon for the user who ran the installation
TARGET_USER="${SUDO_USER}"
if [ -z "$TARGET_USER" ]; then
    # If there's no SUDO_USER, we may be running directly as root
    # Try to find the main user
    TARGET_USER=$(who | head -1 | awk '{print $1}')
fi

if [ -n "$TARGET_USER" ] && [ "$TARGET_USER" != "root" ]; then
    TARGET_UID=$(id -u "$TARGET_USER" 2>/dev/null || true)
    RUNTIME_DIR="/run/user/${TARGET_UID}"
    # `sudo -u` alone doesn't propagate the target user's XDG_RUNTIME_DIR/
    # DBUS_SESSION_BUS_ADDRESS, so `systemctl --user` can't find their session
    # bus and fails silently; set XDG_RUNTIME_DIR explicitly so it can locate
    # /run/user/<uid>/bus. If there's no active session for the user yet (e.g.
    # install ran outside any login session), skip and let them enable it
    # manually instead of failing the whole postinst.
    if [ -n "$TARGET_UID" ] && [ -d "$RUNTIME_DIR" ]; then
        echo "Enabling daemon for user: $TARGET_USER"
        if sudo -u "$TARGET_USER" XDG_RUNTIME_DIR="$RUNTIME_DIR" \
            systemctl --user enable --now hello-daemon.service 2>/dev/null; then
            echo "✓ Daemon enabled and started"
        else
            echo "! Could not enable hello-daemon.service automatically — run 'systemctl --user enable --now hello-daemon.service' as $TARGET_USER"
        fi
    else
        echo "! No active session found for $TARGET_USER — enable the daemon manually with: systemctl --user enable --now hello-daemon.service"
    fi
fi

# PAM configuration is not done here. If libpam-linux-hello is installed,
# its linux-hello-pam-autoconfigure.timer automatically configures sudo and
# screenlock as soon as any user enrolls a face — see
# /usr/share/doc/linux-hello/PAM_MODULE.md, including how to opt out and
# why the SDDM login screen isn't covered by automatic activation.

echo ""
echo "Installation Complete!"
echo ""
echo "Next steps:"
echo "1. Register a face: linux-hello enroll \$(id -u) --context sudo --samples 3"
echo "2. Test with sudo: sudo -v"
echo "3. Or test with screenlock: lock your screen"
echo ""
echo "For help: man linux-hello"
echo ""

exit 0
