# debian/prerm
#!/bin/bash
set -e

if [ "$1" = "remove" ]; then
    systemctl stop linux-hello.service || true
    systemctl disable linux-hello.service || true
fi
