# debian/linux-hello/DEBIAN/postrm
#!/bin/bash
# Post-removal cleanup

set -e

echo "Linux Hello - Cleanup"
echo "===================="

# Stop and disable daemon
systemctl --user stop hello-daemon.service 2>/dev/null || true
systemctl --user disable hello-daemon.service 2>/dev/null || true

echo "✓ Daemon stopped"

# Offer to restore PAM configuration
if compgen -G "/etc/pam.d/sudo.pre-linuxhello-*" > /dev/null; then
    echo ""
    echo "Found PAM backups:"
    find /etc/pam.d -name "*.pre-linuxhello-*" -type f -printf '%T@ %p\n' 2>/dev/null | sort -rn | head -5 | cut -d' ' -f2-
    echo ""
    read -r -p "Restore PAM configuration from backup? (y/n) [y]: " restore_pam
    restore_pam=${restore_pam:-y}
    
    if [ "$restore_pam" = "y" ] || [ "$restore_pam" = "Y" ]; then
        LATEST_SUDO=$(find /etc/pam.d -name "sudo.pre-linuxhello-*" -type f -printf '%T@ %p\n' 2>/dev/null | sort -rn | head -1 | cut -d' ' -f2-)
        if [ -n "$LATEST_SUDO" ]; then
            cp "$LATEST_SUDO" /etc/pam.d/sudo
            echo "✓ Restored /etc/pam.d/sudo from backup"
        fi
    fi
fi

echo ""
echo "Linux Hello uninstalled"
echo ""

exit 0
